# files/harbor/harbor-images-prep.service
[Unit]
Description=Stage k3s air-gap images into /var (first boot)
Before=k3s.service
After=local-fs.target
ConditionPathExists=!/var/lib/rancher/k3s/agent/images/.staged

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/harbor-images-stage

[Install]
WantedBy=multi-user.target
