# files/ensure-mosh-firewall.service
[Unit]
Description=Ensure UFW includes Mosh UDP rules
After=ufw.service
Requires=ufw.service

[Service]
Type=oneshot
ExecStart=/usr/bin/ensure-mosh-firewall

[Install]
WantedBy=multi-user.target
